# etc/cron.d/put-last-post-vdm
5 * * * *  bash /opt/iadvize/library/My/PutLastPostVdm
